# nixos/vms/eve/ansible/roles/php/files/php-fpm-stop.service
[Unit]
Description="Stop php-fpm"

[Service]
ExecStart=/usr/local/bin/stop-php-fpm
