# nixos/vms/eve/ansible/roles/php/files/php-fpm.socket
[Socket]
ListenStream=9000

[Service]
Restart=on-failure

[Install]
WantedBy=sockets.target
