# nixos/vms/eve/ansible/roles/php/files/php-fpm-stop.timer
[Unit]
Description=Stop php-fpm

[Timer]
OnBootSec=2min
OnUnitActiveSec=4min
Persistent=true

[Install]
WantedBy=multi-user.target
